# examples/systemd/llama-server-loopback.service
[Unit]
Description=llama.cpp router on loopback backend port
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/llama
Environment=LLAMA_BIND_HOST=127.0.0.1
Environment=LLAMA_BIND_PORT=8082
Environment=LLAMA_IDLE_UNLOAD_SECONDS=1800
ExecStart=/usr/bin/bash /opt/llama/llama-server.sh
Restart=on-failure
RestartSec=10
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
